# tools/transcript-archive/com.example.transcript-archive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd RENDER TEMPLATE for daily transcript-archive runs.

  This file is not meant to be edited by hand or loaded directly. It is
  rendered by install.sh, which replaces every double-at-sign placeholder
  token below (the ProgramArguments, EnvironmentVariables, schedule, and
  log-path values further down) with a real value and writes the result
  to:
    ~/Library/LaunchAgents/com.skill-issue.transcript-archive.plist
  Run ./install.sh <ARCHIVE_DIR> to do this (see
  skills/claude/transcript-backup/SKILL.md for the guided version).

  IMPORTANT: launchd does NOT inherit your shell environment. The
  TRANSCRIPT_ARCHIVE_DIR / TRANSCRIPT_ARCHIVE_MACHINE_ID variables MUST be
  set here in EnvironmentVariables, not just in your .zshrc / .bashrc --
  install.sh handles this for you.
-->
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.skill-issue.transcript-archive</string>

  <key>ProgramArguments</key>
  <array>
    <string>@@PYTHON@@</string>
    <string>@@BACKUP_PY@@</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TRANSCRIPT_ARCHIVE_DIR</key>
    <string>@@ARCHIVE_DIR@@</string>
    <key>TRANSCRIPT_ARCHIVE_MACHINE_ID</key>
    <string>@@MACHINE_ID@@</string>
  </dict>

  <!-- Daily at the time install.sh was given (default 03:17 local) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>@@HOUR@@</integer>
    <key>Minute</key>
    <integer>@@MINUTE@@</integer>
  </dict>

  <!-- Do NOT run immediately on load -- install.sh runs the first dump in
       the foreground itself; this job only fires at the scheduled time. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>@@LOG_OUT@@</string>

  <key>StandardErrorPath</key>
  <string>@@LOG_ERR@@</string>

</dict>
</plist>
